# templates/thin.init.erb
#!/bin/sh
# file managed by Puppet
### BEGIN INIT INFO
# Provides:          <%= @service %>
# Required-Start:    $local_fs $remote_fs
# Required-Stop:     $local_fs $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: thin initscript for <%= @service %>
# Description:       <%= @service %> init script generated by puppet thin module
### END INIT INFO

PATH=/sbin:/usr/sbin:/bin:/usr/bin:/usr/local/bin
SCRIPT_NAME=/etc/init.d/<%= @service %>
. /lib/lsb/init-functions

<% unless (@force_home.empty?) %>
# Force a HOME environment
# Needed for a puppet master running on thin
export HOME=<%= @force_home %>
<% end %>

run_action() {
        ACTION="$1"
        thin $ACTION --config /etc/thin.d/<%= @name %>.yml
}

case "$1" in
  start)
        run_action start
        ;;
  stop)
        run_action stop
        ;;
  restart|force-reload|reload)
        run_action restart
        ;;
  status)
  <% if (@socket.empty?) %>
        status_of_proc -p /var/run/thin/<%= @name %>.<%= @port %>.pid /usr/local/bin/thin <%= @service %> && exit 0 || exit $?
  <% else %>
        status_of_proc -p /var/run/thin/<%= @name %>.0.pid /usr/local/bin/thin <%= @service %> && exit 0 || exit $?
  <% end %>

        ;;
  *)
        echo "Usage: $SCRIPT_NAME {start|stop|status|restart}" >&2
        exit 3
        ;;
esac

